# systemd/systembackup.service
[Unit]
Description=System daily backup

[Service]
ExecStart=/home/v/bin/systembackup.sh

